# e2e/mac-studio/launchagents/com.boarlabs.e2e.appium.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.boarlabs.e2e.appium</string>

    <!-- Use a login shell so the user's PATH (Homebrew, nvm, etc.) is set up
         before `appium` is resolved. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>exec appium</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Restart only on abnormal exit; don't fight a clean shutdown. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/appium.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/appium.err.log</string>
</dict>
</plist>
